# tools/subscription-cleaner.service.example
[Unit]
Description=Subscription Cleaner Service
After=network.target

[Service]
Type=simple
# TODO: 修改为实际运行用户
User=www-data
Group=www-data
# TODO: 修改为实际部署路径
WorkingDirectory=/opt/subscription-cleaner
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/subscription-cleaner/logs
PrivateTmp=true

[Install]
WantedBy=multi-user.target
