# tools/subscription-cleaner.service.example
[Unit]
Description=Subscription Cleaner Service
After=network.target

[Service]
Type=simple
# 运行用户设置（二选一）：
# 方案1: 使用 root 运行（简单但不推荐）
# User=root
# Group=root

# 方案2: 使用专用用户运行（推荐，需先创建用户）
# 创建用户: sudo useradd --system --user-group --shell /usr/sbin/nologin appuser
# User=appuser
# Group=appuser

WorkingDirectory=/opt/subscription-cleaner
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/subscription-cleaner/logs
PrivateTmp=true

[Install]
WantedBy=multi-user.target
